# rasterblocks/package-rasterblocks/DEBIAN/postinst
#!/bin/bash

INSTALL_DIR=/usr/lib/rasterblocks

systemctl enable "$INSTALL_DIR/rasterblocks.service"
systemctl start rasterblocks.service
